# shared/n1.upstart.conf
description "etcd service registry"

start on started
stop on shutdown

exec /usr/bin/consul agent -bind 172.20.20.10 -server -bootstrap -data-dir /tmp/consul