# shared/n3.upstart.conf
description "etcd service registry"

start on started
stop on shutdown

exec /usr/bin/consul agent -bind 172.20.20.12 -server -data-dir /tmp/consul